# env/image/mhtweb/files/etc/supervisord.d/crond.ini
[program:crond]
process_name= crond
directory= /etc/
command= /usr/sbin/crond start
startsecs= 0
autorestart= false
