# env/image/mhtweb/files/etc/supervisord.d/httpd.ini
[program:httpd]
command=/usr/sbin/httpd -D FOREGROUND
autostart=true
autorestart=true

